# laradock/stacks/alpine_nginx/nginx/etc/service/nginx/run
#!/bin/sh

exec 2>&1
source /env

mkdir -p /run/nginx
exec /usr/sbin/nginx -g 'daemon off;'
