# laradock/stacks/alpine_apache/apache/etc/service/apache/run
#!/bin/sh

exec 2>&1
source /env

exec /usr/sbin/httpd -D FOREGROUND
