# laradock/stacks/alpine_nginx/nginx/etc/service/php-fpm/run
#!/bin/sh

exec 2>&1
source /env

exec php-fpm --nodaemonize
